# config/autostart/fcitx5.desktop
[Desktop Entry]
Type=Application
Name=Fcitx5
Exec=fcitx5 -d
StartupNotify=false
Terminal=false
